# games/elda-go/src/elda-go/x52mfd.service.in
# Simple systemd scenario for x52mfd
# - put it in @HOME@/.config/systemd/user/
# - run: systemd --user reload-daemon
# - use: systemctl --user <status|start|stop> x52mfd
###########################################################

[Unit]
Description = x52 joystick game support for Elite Dangerous

[Service]
Environment="ELDA_DEBUG=1"
Type = simple
StandardOutput = journal
ExecStart = @HOME@/.local/bin/x52mfd elda-go

[Install]
WantedBy=default.target

